# archiso-ffx1/airootfs/etc/systemd/system/set-factory-boot.service
[Unit]
Description=Set systemd-boot to factory_reset.conf on next boot

[Service]
Type=simple
ExecStart=bash -c 'bootctl set-oneshot factory_reset.conf && sync && sleep 5 && systemctl reboot'