# archiso-ff1/airootfs/etc/systemd/system/earlyoom.service.d/override.conf
[Service]
ExecStart=
ExecStart=earlyoom -m 10 -s 100 -r 0 -g -N /usr/local/bin/ff1-earlyoom-hook.sh --prefer '^(chromium|cage)' --avoid '^(feral-|systemd)' -n --syslog
