# archiso-ff1/airootfs/etc/systemd/system/earlyoom.service.d/override.conf
[Service]
DynamicUser=false
ExecStart=
ExecStart=earlyoom -m 10 -s 100 -r 0 -g -N '/usr/local/bin/earlyoom-notify.sh' --prefer '^(chromium|cage)' --avoid '^(feral-|systemd)' -n --syslog
