# archiso-ff1/airootfs/etc/systemd/system/set-factory-boot.service
[Unit]
Description=Set systemd-boot to factory_reset.conf on next boot

[Service]
Type=simple
ExecStart=/root/scripts/factory_reset.sh