# archiso-ff1/airootfs/etc/NetworkManager/dispatcher.d/09-timezone
#!/bin/bash

# NetworkManager dispatcher script for automatic timezone synchronization
# Place this in /etc/NetworkManager/dispatcher.d/
# Make executable with: chmod +x /etc/NetworkManager/dispatcher.d/99-timezone-sync

TIMEOUT=5

# Service definitions: URL|JQ_FILTER
# This list is ordered by priority
SERVICES='
https://ipapi.co/json/|.timezone
https://ipinfo.io/json|.timezone
https://geoip.chrisdown.name/|.location.time_zone
https://worldtimeapi.org/api/ip/|.timezone
https://ipwho.is/|.timezone.id
https://reallyfreegeoip.org/json/|.time_zone
https://api.ipbase.com/v1/json/|.time_zone
'

# Extract timezone from JSON
extract_tz() {
    local json="$1"
    local filter="$2"
    
    echo "$json" | jq -r "$filter" 2>/dev/null | grep -v "^null$"
}

validate_timezone() {
    local tz="$1"
    if timedatectl list-timezones 2>/dev/null | grep -qx "$tz"; then
        return 0
    fi
    return 1
}

# Query single service
query_service() {
    local url="$1"
    local filter="$2"
    local response
    local tz
    
    response=$(curl -sf --max-time "$TIMEOUT" "$url" 2>/dev/null)
    if [ $? -eq 0 ] && [ -n "$response" ]; then
        tz=$(extract_tz "$response" "$filter")
        # Validate timezone format
        if [ -n "$tz" ] && validate_timezone "$tz"; then
            echo "$tz"
            return 0
        fi
    fi
    return 1
}

# Get timezone (first successful response)
get_timezone_first() {
    while IFS='|' read -r url filter; do
        [ -z "$url" ] && continue
        
        logger -t "nm-timezone" "Trying service: $url"
        tz=$(query_service "$url" "$filter")
        
        if [ -n "$tz" ]; then
            echo "$tz"
            return 0
        fi
    done <<< "$SERVICES"
    return 1
}

# Main dispatcher logic
case "$2" in
    up)
        timezone=$(get_timezone_first)
        
        # Set timezone if valid
        if [ -n "$timezone" ]; then
            logger -t "nm-timezone" "Setting timezone to: $timezone"
            timedatectl set-timezone "$timezone" 2>&1 | logger -t "nm-timezone"
        else
            logger -t "nm-timezone" "Failed to determine timezone"
            exit 1
        fi
        ;;
esac