# archiso-ff1/airootfs/etc/NetworkManager/dispatcher.d/09-timezone
#!/bin/bash

# NetworkManager dispatcher script for automatic timezone synchronization
# Path: /etc/NetworkManager/dispatcher.d/99-timezone-sync
# Permissions: chmod +x /etc/NetworkManager/dispatcher.d/99-timezone-sync
# Dependencies: curl, jq

TIMEOUT=5
MAX_RETRIES=3

# Service definitions: URL|JQ_FILTER
SERVICES='
https://ipapi.co/json/|.timezone
https://ipinfo.io/json|.timezone
https://geoip.chrisdown.name/|.location.time_zone
https://worldtimeapi.org/api/ip/|.timezone
https://ipwho.is/|.timezone.id
https://reallyfreegeoip.org/json/|.time_zone
https://api.ipbase.com/v1/json/|.time_zone
'

extract_tz() {
    local json="$1"
    local filter="$2"
    echo "$json" | jq -r "$filter" 2>/dev/null | grep -v "^null$"
}

validate_timezone() {
    local tz="$1"
    if timedatectl list-timezones 2>/dev/null | grep -qx "$tz"; then
        return 0
    fi
    return 1
}

query_service() {
    local url="$1"
    local filter="$2"
    local response
    local tz
    
    if response=$(curl -sf --max-time "$TIMEOUT" "$url" 2>/dev/null) && [ -n "$response" ]; then
        tz=$(extract_tz "$response" "$filter")
        if [ -n "$tz" ] && validate_timezone "$tz"; then
            echo "$tz"
            return 0
        fi
    fi
    return 1
}

get_timezone_first() {
    local attempt=0

    while [ $attempt -lt $MAX_RETRIES ]; do
        local url filter tz
        while IFS='|' read -r url filter; do
            [[ -z "$url" || "$url" == \#* ]] && continue
            
            logger -t "nm-timezone" "Querying: $url (attempt $((attempt + 1))/$MAX_RETRIES)"
            tz=$(query_service "$url" "$filter")
            
            if [ -n "$tz" ]; then
                echo "$tz"
                return 0
            fi
        done <<< "$SERVICES"
        
        ((attempt++))
        if [ $attempt -lt $MAX_RETRIES ]; then
            logger -t "nm-timezone" "All services failed, retrying in 2 seconds..."
            sleep 2
        fi
    done
    
    logger -t "nm-timezone" "Failed to fetch valid timezone after $MAX_RETRIES attempts"
    return 1
}

# Main Logic
case "$2" in
    up)
        # Check if manual override flag exists
        if [ -f /etc/timezone-manual-override ]; then
            logger -t "nm-timezone" "Manual override flag detected (/etc/timezone-manual-override). Skipping auto-sync."
            exit 0
        fi
        
        timezone=$(get_timezone_first)
        
        if [ -n "$timezone" ]; then
            current_tz=$(timedatectl show --property=Timezone --value 2>/dev/null)
            if [ "$timezone" != "$current_tz" ]; then
                logger -t "nm-timezone" "Changing timezone from '$current_tz' to '$timezone'"
                if timedatectl set-timezone "$timezone" 2>&1 | logger -t "nm-timezone"; then
                    logger -t "nm-timezone" "Timezone successfully updated."
                fi
            else
                logger -t "nm-timezone" "Timezone is already correct ($timezone). No change needed."
            fi
        else
            logger -t "nm-timezone" "Failed to fetch valid timezone from all services."
            exit 1
        fi
        ;;
esac
